# deploy/twy-docs.service
[Unit]
Description=TWY User Guide (docs.tiffanywoodyoga.com)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/twy/overview
EnvironmentFile=/root/twy/overview/.env.docs
ExecStart=/usr/bin/python3 /root/twy/overview/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
